# deploy/systemd/turbalance-machine-benchmark.timer
[Unit]
Description=Run Turbalance machine benchmarks periodically

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
AccuracySec=1min
Persistent=true
RandomizedDelaySec=90s

[Install]
WantedBy=timers.target
